# ros2_ws/src/scoutmini_streaming/systemd/scoutmini-webrtc.service
[Unit]
Description=ScoutMini ZED WebRTC gateway
After=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
Environment=SCOUT_REPO_ROOT=%h/repos/ScoutMini
ExecStart=%h/repos/ScoutMini/ros2_ws/install/scoutmini_streaming/share/scoutmini_streaming/scripts/start_zed_webrtc_stack.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
